# src/skel/etc/supervisor.d/program-logstash.conf
; -*- ini -*-
[program:logstash]
command = /opt/logstash/bin/launch.sh
