# src/skel/etc/supervisor.d/program-incrond.conf
; -*- ini -*-
[program:incrond]
command = /usr/sbin/incrond --foreground --config=/etc/incron.conf
