# src/supervisord.conf
[program:logstash]
command = /usr/local/bin/launch-logstash.sh
